# crontab
* * * * * /var/www/site.com/test1.php

# Some comment
* * * * * /var/www/site.com/test2.php
